# deploy/nginx/redis.ini
; ================================
;  redis supervisor
; ================================

[program:redis]
command=/usr/bin/redis-server /home/transtats/transtats/deploy/nginx/redis.conf
directory=/home/transtats/transtats

;user=root
numprocs=1
stdout_logfile=/home/transtats/transtats/transtats/logs/celery/redis.log
stderr_logfile=/home/transtats/transtats/transtats/logs/celery/redis_err.log
autostart=true
autorestart=true
startsecs=10

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998
